# conf/templates/billing-worker.conf.erb
description "celery server"
start on runlevel [2345]
script
su --login  <%= @username%> --session-command "cd /var/local/<%= @username%>/billing/;celery -A core.extraction.task worker --loglevel=info"
end script
console output
respawn
respawn limit 2 10
